# project/runnix/target_skeleton/etc/inittab
console::sysinit:/etc/rc install

::ctrlaltdel:/sbin/reboot
::shutdown:/etc/rc shutdown
