# project/astlinux/board/wrap/target_skeleton/etc/inittab
# 
# /etc/inittab
#
# $Id: inittab 1681 2004-09-01 18:12:49Z  $
#

console::sysinit:/etc/rc sysinit
console::respawn:/sbin/getty -L 38400 /dev/ttyS0 vt100

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/rc shutdown
